# netmon.service
[Unit]
Description=Network connectivity moniting

[Service]
Type=simple
ExecStart=/path/to/network-mon

[Install]
WantedBy=multi-user.target